# supervisord.conf
[supervisord]
nodaemon=true

[program:tomcat]
command=/root/webserver/tomcat/bin/startup.sh
